# lsio-mods/emby/root-fs/etc/s6-overlay/s6-rc.d/init-mod-sqlite3-preflight/up
/etc/s6-overlay/s6-rc.d/init-mod-sqlite3-preflight/run
