# lsio-mods/plex/root-fs/etc/s6-overlay/s6-rc.d/init-mod-sqlite3-preflight/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
# shellcheck source=/dev/null
set -u
phase_name="01-preflight"
. /opt/sqlite3-lsio-mod/lib/logging.sh
. /opt/sqlite3-lsio-mod/lib/sha.sh
. /opt/sqlite3-lsio-mod/lib/arch.sh
. /opt/sqlite3-lsio-mod/lib/manifest-parser.sh
. /opt/sqlite3-lsio-mod/lib/selector.sh
sqlite3_mod_init_logging
baked_pins="/opt/sqlite3-lsio-mod/baked-pins.txt"
if [ ! -f "$baked_pins" ]; then log_fatal "mod=plex phase=${phase_name} event=missing-baked-pins"; fi
# Phase 01 command checks are a best-effort early warning; later phases re-check their mutation prerequisites.
if ! require_all_or_warn awk chmod chown cp grep mkdir mktemp mv rm sed sha256sum stat tr uname; then exit 0; fi
validate_baked_pins_schema "$baked_pins" || log_fatal "mod=plex phase=${phase_name} event=malformed-baked-pins"
selected_arch="$(resolve_arch)" || { log_warn "mod=plex phase=${phase_name} event=unsupported-arch machine=$(uname -m)"; exit 0; }
server_id="$(select_supported_server "plex" "$selected_arch" "$baked_pins")"
[ -n "$server_id" ] || exit 0
if unsupported_row="$(selected_unsupported_row "plex" "$selected_arch" "$baked_pins" "$server_id")"; then
  IFS='|' read -r unsupported_compat_group unsupported_reason <<< "$unsupported_row"
  log_warn "mod=plex phase=${phase_name} event=unsupported-server arch=${selected_arch} server_id=${server_id} compat_group=${unsupported_compat_group} reason=${unsupported_reason}"
  exit 0
fi
artifact_row="$(selected_artifact_row "plex" "$selected_arch" "$baked_pins" "$server_id")" || log_fatal "mod=plex phase=${phase_name} event=missing-artifact-row arch=${selected_arch} server_id=${server_id}"
IFS='|' read -r compat_group artifact_relpath target_path artifact_sha256 <<< "$artifact_row"
: "$artifact_relpath" "$artifact_sha256"
[ -n "$compat_group" ] || log_fatal "mod=plex phase=${phase_name} event=invalid-artifact-row arch=${selected_arch} server_id=${server_id}"
[ -n "$target_path" ] || log_fatal "mod=plex phase=${phase_name} event=invalid-artifact-row arch=${selected_arch} server_id=${server_id}"
if [ ! -f "$target_path" ]; then
  log_warn "mod=plex phase=${phase_name} event=missing-target arch=${selected_arch} server_id=${server_id} compat_group=${compat_group} target_path=${target_path}"
  exit 0
fi
log_info "mod=plex phase=${phase_name} event=selected arch=${selected_arch} server_id=${server_id} compat_group=${compat_group} target_path=${target_path}"
