# lsio-mods/plex/root-fs/etc/s6-overlay/s6-rc.d/init-mod-sqlite3-poolpatch/run
#!/usr/bin/with-contenv bash
set -u
phase_name="04-pool-patch"
. /opt/sqlite3-lsio-mod/lib/logging.sh
. /opt/sqlite3-lsio-mod/lib/sha.sh
. /opt/sqlite3-lsio-mod/lib/arch.sh
. /opt/sqlite3-lsio-mod/lib/atomic-write.sh
. /opt/sqlite3-lsio-mod/lib/plex-pool-patch.sh
sqlite3_mod_init_logging
case "$(uname -m)" in
  aarch64|arm64) log_warn "mod=plex phase=${phase_name} event=pool-patch-deferred arch=linux-arm64"; exit 0 ;;
esac
if ! require_all_or_warn awk chmod chown cp dd grep mktemp mv od printf rm sha256sum stat tr uname; then exit 0; fi
baked_pins="/opt/sqlite3-lsio-mod/baked-pins.txt"
target_path="/usr/lib/plexmediaserver/lib/libsqlite3.so"
selected_arch="$(resolve_arch)" || exit 0
current_sha="$(awk -F'|' -v arch="$selected_arch" -v target="$target_path" '$1=="current" && $3=="plex" && $4==arch && $6==target { print $7; found=1 } END { if (!found) exit 1 }' "$baked_pins")" || exit 0
[ -f "$target_path" ] && [ "$(sha256_of "$target_path")" = "$current_sha" ] || { log_warn "mod=plex phase=${phase_name} event=skip-sqlite-not-current"; exit 0; }
pms="/usr/lib/plexmediaserver/Plex Media Server"
scanner="/usr/lib/plexmediaserver/Plex Media Scanner"
pms_sha="$(awk -F'|' -v arch="$selected_arch" -v path="$pms" '$1=="pool-pre" && $3=="plex" && $4==arch && $5==path { print $6; found=1 } END { if (!found) exit 1 }' "$baked_pins")" || { log_warn "mod=plex phase=${phase_name} event=missing-pool-row binary=pms"; exit 0; }
scanner_sha="$(awk -F'|' -v arch="$selected_arch" -v path="$scanner" '$1=="pool-pre" && $3=="plex" && $4==arch && $5==path { print $6; found=1 } END { if (!found) exit 1 }' "$baked_pins")" || { log_warn "mod=plex phase=${phase_name} event=missing-pool-row binary=scanner"; exit 0; }
plex_pool_patch_apply_binary "$pms" "$pms_sha" \
  "0xae4a17|11422231|11422232|be140000004c89ff41b800000000506a|be100000004c89ff41b800000000506a" || true
plex_pool_patch_apply_binary "$scanner" "$scanner_sha" \
  "0x36f06d|3600493|3600494|be140000004c89ff41b800000000506a|be100000004c89ff41b800000000506a" \
  "0x38b37b|3715963|3715964|be140000004889df41b80000000068d0|be100000004889df41b80000000068d0" || true
log_info "mod=plex phase=${phase_name} event=complete"
