# lsio-mods/plex/root-fs/etc/s6-overlay/s6-rc.d/init-mod-sqlite3-poolpatch/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
# shellcheck source=/dev/null
set -u
phase_name="04-pool-patch"
. /opt/sqlite3-lsio-mod/lib/logging.sh
. /opt/sqlite3-lsio-mod/lib/sha.sh
. /opt/sqlite3-lsio-mod/lib/arch.sh
. /opt/sqlite3-lsio-mod/lib/manifest-parser.sh
. /opt/sqlite3-lsio-mod/lib/selector.sh
. /opt/sqlite3-lsio-mod/lib/atomic-write.sh
. /opt/sqlite3-lsio-mod/lib/plex-pool-patch.sh
sqlite3_mod_init_logging
if ! require_all_or_warn awk chmod chown cp dd grep mktemp mv od printf rm sha256sum stat tr uname; then exit 0; fi
baked_pins="/opt/sqlite3-lsio-mod/baked-pins.txt"
selected_arch="$(resolve_arch)" || exit 0
server_id="$(select_supported_server "plex" "$selected_arch" "$baked_pins")"
[ -n "$server_id" ] || exit 0
if unsupported_row="$(selected_unsupported_row "plex" "$selected_arch" "$baked_pins" "$server_id")"; then
  IFS='|' read -r unsupported_compat_group unsupported_reason <<< "$unsupported_row"
  log_warn "mod=plex phase=${phase_name} event=unsupported-server arch=${selected_arch} server_id=${server_id} compat_group=${unsupported_compat_group} reason=${unsupported_reason}"
  exit 0
fi
if ! selected_artifact_is_current "plex" "$selected_arch" "$baked_pins" "$server_id"; then
  log_warn "mod=plex phase=${phase_name} event=skip-sqlite-not-current"
  exit 0
fi
pms_detector_row="$(selected_pristine_detector_row "plex" "$selected_arch" "$baked_pins" "$server_id" "plex_pms")" || { log_warn "mod=plex phase=${phase_name} event=missing-pool-row binary=pms"; exit 0; }
IFS='|' read -r pms_detector_path pms_sha <<< "$pms_detector_row"
if [ -z "$pms_detector_path" ] || [ -z "$pms_sha" ]; then
  log_warn "mod=plex phase=${phase_name} event=missing-pool-row binary=pms"
  exit 0
fi
scanner_detector_row="$(selected_pristine_detector_row "plex" "$selected_arch" "$baked_pins" "$server_id" "plex_scanner")" || { log_warn "mod=plex phase=${phase_name} event=missing-pool-row binary=scanner"; exit 0; }
IFS='|' read -r scanner_detector_path scanner_sha <<< "$scanner_detector_row"
if [ -z "$scanner_detector_path" ] || [ -z "$scanner_sha" ]; then
  log_warn "mod=plex phase=${phase_name} event=missing-pool-row binary=scanner"
  exit 0
fi
pms_site_rows="$(selected_pool_site_rows "plex" "$selected_arch" "$baked_pins" "$server_id" "$pms_detector_path")" || { log_warn "mod=plex phase=${phase_name} event=missing-pool-row binary=pms"; exit 0; }
scanner_site_rows="$(selected_pool_site_rows "plex" "$selected_arch" "$baked_pins" "$server_id" "$scanner_detector_path")" || { log_warn "mod=plex phase=${phase_name} event=missing-pool-row binary=scanner"; exit 0; }
pms_sites=()
while IFS= read -r site_row || [ -n "$site_row" ]; do
  [ -n "$site_row" ] || continue
  pms_sites+=("$site_row")
done <<EOF_PMS_SITES
$pms_site_rows
EOF_PMS_SITES
scanner_sites=()
while IFS= read -r site_row || [ -n "$site_row" ]; do
  [ -n "$site_row" ] || continue
  scanner_sites+=("$site_row")
done <<EOF_SCANNER_SITES
$scanner_site_rows
EOF_SCANNER_SITES
plex_pool_patch_apply_binary "$pms_detector_path" "$pms_sha" "${pms_sites[@]}" || true
plex_pool_patch_apply_binary "$scanner_detector_path" "$scanner_sha" "${scanner_sites[@]}" || true
log_info "mod=plex phase=${phase_name} event=complete"
